# auto-analisa-web/deploy/systemd/autoanalisa-macro.timer
[Unit]
Description=Run Macro Daily Generator at 00:05 and 18:05 WIB

[Timer]
# Catatan: Systemd menggunakan timezone host. Pastikan host diset ke Asia/Jakarta
# atau sesuaikan jadwal OnCalendar berikut ke UTC bila host bukan WIB.
OnCalendar=*-*-* 00:05:00
OnCalendar=*-*-* 18:05:00
Persistent=true
Unit=autoanalisa-macro.service

[Install]
WantedBy=timers.target
